# services/hadoop/supervisord.conf
[supervisord]
nodaemon=true

[program:setup]
command = /tmp/setup.sh
startsecs = 0
autorestart = false
startretries = 1

[program:namenode]
command=/opt/hadoop/bin/hdfs --config /opt/hadoop/etc/hadoop namenode
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:datanode]
command=/opt/hadoop/bin/hdfs --config /opt/hadoop/etc/hadoop datanode
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:resourcemanager]
command=/opt/hadoop/bin/yarn --config /opt/hadoop/etc/hadoop resourcemanager
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nodemanager]
command=/opt/hadoop/bin/yarn --config /opt/hadoop/etc/hadoop nodemanager
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

